# shinken-thruk/files/thruk.conf
[program:httpd]
command=/usr/sbin/httpd -f /etc/httpd/conf/httpd.conf -DNO_DETACH -DFOREGROUND
autostart=true
logfile=/var/log/supervisor-httpd.log    