# files/supervisor/mongodb.conf
[program:mongodb]
command=/usr/bin/mongod
autostart=true
log_stdout=true
log_stderr=true
logfile=/var/log/supervisor-mongodb.log    
