# files/supervisor/shinken-poller.conf
[program:shinken-poller]
command=python /usr/local/shinken/bin/shinken-poller -c /etc/shinken/pollerd.ini
autostart=true
logfile=/var/log/supervisor-shinken-poller.log    
