# files/supervisor/shinken-poller.conf
[program:shinken-poller]
command=/usr/bin/shinken-poller -c /etc/shinken/daemons/pollerd.ini
autostart=true
logfile=/var/log/supervisor-shinken-poller.log    
